# scripts/feed-br-history.service
[Unit]
Description=Teledigest — feed historical BR messages to mining pipeline
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Container name is volatile (Dokploy rebuilds change suffix). Find it by
# the image command at exec time so this works across redeploys.
ExecStart=/bin/bash -c '\
  # Dokploy / Docker Swarm gives containers names like:
  #   bots-grab-rlcj1h.1.h816na530nwfpt1q9r6nc34d2
  # Image names also have the project prefix. Adjust the grep below if your
  # deployment uses a different project name.
  CID=$(docker ps --format "{{.ID}} {{.Names}} {{.Image}}" \
    | grep -iE "bots-grab|teledigest" | awk "{print \\$1}" | head -1); \
  if [ -z "$CID" ]; then \
    echo "feed-br-history: no teledigest container found" >&2; \
    echo "feed-br-history: docker ps shows:" >&2; \
    docker ps --format "  {{.ID}} {{.Names}} {{.Image}}" >&2; \
    exit 1; \
  fi; \
  echo "feed-br-history: using container $CID"; \
  docker exec "$CID" python -m teledigest.scripts.feed_br_history --country br --batch 15'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
